# profile.d/android.sh
export ANDROID_HOME=$HOME/Library/Android/sdk
export JAVA_HOME=/usr/local/opt/openjdk@21
export PATH=$JAVA_HOME/bin:$ANDROID_HOME/emulator:$ANDROID_HOME/platform-tools:$ANDROID_HOME/tools:$ANDROID_HOME/tools/bin:$PATH

