# profile.d/flutter.sh
export PATH=$HOME/Library/flutter/bin:$HOME/.pub-cache/bin:$PATH
