# profile.d/node.sh
export PATH=/usr/local/opt/node@20/bin:$PATH
