# profile.d/virtualenv.sh
#Where to work on virtualenvs
export WORKON_HOME=$HOME/.virtualenv

#Where projects will reside
export PROJECT_HOME=$HOME/virtualenv-projects

#Location of virtualenvwrapper.sh
source /usr/local/bin/virtualenvwrapper.sh
