# rules.d/91-garmin.rules
SUBSYSTEM=="usb", ATTRS{idVendor}=="091e", ATTRS{idProduct}=="24f0",RUN+="/home/slmolloy/projects/garmin/runner.sh"
